# influxdb3/rootfs/etc/s6-overlay/s6-rc.d/init-influxdb3/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: InfluxDB3
# Configures InfluxDB3
# ==============================================================================
declare INFLUXDB3_ENTERPRISE_LICENSE_EMAIL
declare INFLUXDB3_ENTERPRISE_LICENSE_TYPE
declare INFLUXDB3_ENTERPRISE_CLUSTER_ID
declare INFLUXDB3_NODE_IDENTIFIER_PREFIX
declare INFLUXDB3_OBJECT_STORE
declare INFLUXDB3_DB_DIR
declare LOG_FILTER

INFLUXDB3_ENTERPRISE_LICENSE_EMAIL=$(bashio::config 'license_email')

INFLUXDB3_ENTERPRISE_LICENSE_TYPE=home
INFLUXDB3_ENTERPRISE_CLUSTER_ID=cluster0
INFLUXDB3_NODE_IDENTIFIER_PREFIX=my-node
INFLUXDB3_OBJECT_STORE=file
INFLUXDB3_DB_DIR=/data/influxdb3
LOG_FILTER=info