# influxdb3/rootfs/etc/s6-overlay/s6-rc.d/influxdb3/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: InfluxDB3
# Runs the InfluxDB3 Server
# ==============================================================================
export INFLUXDB3_ENTERPRISE_LICENSE_EMAIL=$(bashio::config 'license_email')
export INFLUXDB3_ENTERPRISE_LICENSE_TYPE=home
export INFLUXDB3_ENTERPRISE_CLUSTER_ID=cluster0
export INFLUXDB3_NODE_IDENTIFIER_PREFIX=my-node
export INFLUXDB3_OBJECT_STORE=file
export INFLUXDB3_DB_DIR=/data/influxdb3
export LOG_FILTER=$(bashio::config 'influxd_log_level')

for envvar in $(bashio::config 'envvars|keys'); do
    name=$(bashio::config "envvars[${envvar}].name")
    value=$(bashio::config "envvars[${envvar}].value")
    bashio::log.debug "Setting Env Variable ${name} to ${value}"
    export "${name}=${value}"
done

bashio::log.info 'Starting InfluxDB v3.x...'

# Run InfluxDB
exec influxdb3 serve --plugin-dir /data/influx-plugins